# supervisord/delayed_job.conf
[program:rails]
user=openproject
directory=/var/www/openproject
command=/var/www/openproject/docker/scripts/env.sh bundle exec script/delayed_job run --exit-on-complete
autorestart=true
